# scripts/drone-relay.service
[Unit]
Description=Drone Bridge UDP Relay Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/drone-relay
ExecStart=/usr/bin/python3 /opt/drone-relay/udp_relay_vps.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
